# ops/systemd/thegame-backend.service
[Unit]
Description=TheGame API (Fastify)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/dev/thegame/backend
User=dev
Group=dev

Environment=NODE_ENV=production
Environment=PATH=__NODE_BIN__:/usr/bin:/bin
# Put environment variables in /etc/thegame/thegame-backend.env (required for DATABASE_URL, JWT secrets, etc)
EnvironmentFile=-/etc/thegame/thegame-backend.env

ExecStart=__NODE_BIN__/node dist/server.js
Restart=always
RestartSec=3

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
