# ops/systemd-user/thegame-backend.service
[Unit]
Description=TheGame API (Fastify) - user service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/dev/thegame/backend

Environment=NODE_ENV=production
Environment=PATH=__NODE_BIN__:/usr/bin:/bin
# Put environment variables in ~/.config/thegame/thegame-backend.env (required for DATABASE_URL, JWT secrets, etc)
EnvironmentFile=%h/.config/thegame/thegame-backend.env

ExecStart=__NODE_BIN__/node dist/server.js
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
