# ops/systemd-user/thegame-monitor.service
[Unit]
Description=TheGame monitor (user)

[Service]
Type=oneshot
ExecStart=/home/dev/thegame/ops/scripts/check-services-user.sh
